# packages/kadi/kadi-0.9.0-py3-none-any.whl/kadi/cli/templates/kadi-celery.service
[Unit]
Description=Kadi Celery service
After=network.target
After=postgresql.service
After=elasticsearch.service
After=redis-server.service

[Service]
ExecStart={{ kadi_bin }} celery worker --loglevel=INFO --logfile=/var/log/celery/%n%I.log --pidfile=/var/run/celery/%n.pid
Type=simple
Restart=always
KillMode=mixed
TimeoutStopSec=3600

UMask=027
User={{ uid }}
Group={{ gid }}

Environment=KADI_CONFIG_FILE={{ kadi_config }}

[Install]
WantedBy=multi-user.target
