# packages/clustermgr4/clustermgr4-4.2.post5.tar.gz/clustermgr4-4.2.post5/clustermgr/templates/stunnel/stunnel.service
# Stunnel Service file for CentOS 7

[Unit]
Description=SSL tunnel for network daemons
After=network.target
After=syslog.target

[Install]
WantedBy=multi-user.target
Alias=stunnel.target

[Service]
Type=forking
ExecStart=/usr/bin/stunnel /etc/stunnel/stunnel.conf
TimeoutSec=600
Restart=always
PrivateTmp=false
