# packages/botd/botd-33.tar.gz/botd-33/botd.service
[Unit]
Description=BOTD - 24/7 channel daemon
After=multi-user.target

[Service]
DynamicUser=True
StateDirectory=botd
LogsDirectory=botd
CacheDirectory=botd
ExecStart=/usr/local/bin/botd
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
