# packages/shyft/shyft-4.23.3-cp38-cp38-manylinux1_x86_64.whl/shyft-4.23.3.data/purelib/shyft/energy_market/service/install/dstm.service
[Unit]
Description=DistributedShortTermModelsHotService

[Service]
Type=simple
ExecStart=/usr/local/bin/dstm.sh ${web_ip} ${port} ${web_port} ${web_root} ${log_root} ${conda_env} ${dtss_config}
Restart=on-failure

[Install]
WantedBy=multi-user.target
