# packages/understory/understory-0.0.2-py3-none-any.whl/understory/web/framework/templates/supervisor.conf
$def with (root, root_hash, user, env, app, bot)
[program:$(app)_app_$root_hash]
environment = WEBCFG=$root/etc/web.conf,KVDB=$root/run/redis.sock,PYTHONUNBUFFERED=1
command = $env/src/pkg/pkg/runinenv $env web serve $app
directory = $root
user = $user
autostart = true
stopsignal = INT

[program:$(app)_bot_$root_hash]
environment = PATH="/home/gaea/detritus/bin:%(ENV_PATH)s",WEBCFG=$root/etc/web.conf,KVDB=$root/run/redis.sock,PYTHONUNBUFFERED=1
command = $env/src/pkg/pkg/runinenv $env $bot jobqueue
directory = $root
user = $user
autostart = true
stopsignal = INT

[program:$(app)_redis_$root_hash]
command = redis-server etc/redis.conf
directory = $root
user = $user
autostart = true
stopsignal = INT

[program:$(app)_nginx_$root_hash]
command = /home/gaea/detritus/nginx/sbin/nginx
directory = $root
user = root
autostart = true
stopsignal = INT

[program:$(app)_tor_$root_hash]
command = /home/gaea/detritus/bin/tor
directory = $root
user = $user
autostart = true
stopsignal = INT

[group:$(app)_$root_hash]
programs = $(app)_app_$root_hash,$(app)_bot_$root_hash,$(app)_redis_$root_hash,$(app)_nginx_$root_hash,$(app)_tor_$root_hash
