# packages/shyft/shyft-4.23.3-cp37-cp37m-manylinux1_x86_64.whl/shyft-4.23.3.data/purelib/shyft/energy_market/service/install/stm_task.service
[Unit]
Description=ShortTermModelTaskService

[Service]
Type=simple
ExecStart=/usr/local/bin/stm_task.sh ${web_ip} ${port} ${web_port} ${root_dir} ${model_root} ${conda_env} ${log_file}
Restart=on-failure

[Install]
WantedBy=multi-user.target
